# tftdisplay.service
[Unit]
Description=TFT Display for Audio Service
After=network.target

[Service]
Type=simple
User=admin
Group=admin
ExecStart=/home/admin/tft-display/tft-display.py
WorkingDirectory=/home/admin/tft-display
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
